# bash/Ubuntu/update-motd.d/10-help-text
#!/bin/sh
printf "\n"
printf " * Documentation:  https://koliasa.com\n"
printf " * Management:     https://fpv.sale\n"
